# ftp.cron
30 0-23 * * * /media/sf_vmshared/ftpdownloader/ftpdownloader.sh
