# ftp.cron
*/10 * * * * /media/sf_vmshared/ftpdownloader/ftpdownloader.sh
